# deploy/qc.timer
[Unit]
Description=Healthcare Data QC Pipeline — poll every 15 minutes

[Timer]
OnBootSec=5min
OnUnitInactiveSec=15min
AccuracySec=1min

[Install]
WantedBy=timers.target
